# com.sentryclaw.yt-worker-audio.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.sentryclaw.yt-worker-audio</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/sentryclaw/Projects/youtube-transcriber/scripts/start_worker.sh</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>CELERY_HOSTNAME</key>
    <string>native-audio-worker@%h</string>
    <key>CELERY_QUEUES</key>
    <string>audio</string>
  </dict>

  <key>WorkingDirectory</key>
  <string>/Users/sentryclaw/Projects/youtube-transcriber</string>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/yt-worker/yt-worker-audio.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/yt-worker/yt-worker-audio.log</string>
</dict>
</plist>
