# com.sentryclaw.yt-chatbot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sentryclaw.yt-chatbot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/sentryclaw/Projects/youtube-transcriber/scripts/start_telegram_bot.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/sentryclaw/Projects/youtube-transcriber</string>

    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/yt-chatbot/yt-chatbot.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/yt-chatbot/yt-chatbot.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
